# src/systemd/gdevrefresh.service.in
[Unit]
Description=Logitech Gaming Devices Refresh
Requires=gdevd.service
After=suspend.target
After=hibernate.target
After=hybrid-sleep.target

[Service]
Type=oneshot
ExecStart=$$PREFIX$$/bin/gdevctl refresh

[Install]
WantedBy=suspend.target
WantedBy=hibernate.target
WantedBy=hybrid-sleep.target